# devices/common/diy/package/base-files/files/etc/crontabs/root
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
0 */1 * * * /etc/hosts.sh
0 1 * * * /etc/updatew.sh
