# devices/common/diy/package/base-files/files/etc/crontabs/root
0 */3 * * * /etc/hosts.sh
0 1 * * * /etc/updatew.sh
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
0 3 * * * /etc/mosdns/script/update.sh
0 4 * * * /etc/mosdns/script/flush_cache.sh
